# recipes-example/init_dhcpcd/files/init_dhcpcd_backup_250623
### BEGIN INIT INFO
# Provides:          dhcpcd
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start dhcpcd at boot time
# Description:       Enable DHCP client daemon
### END INIT INFO

case "$1" in
  start)
    /sbin/dhcpcd
    sleep 1
    /usr/bin/debug-app ethernet init
    sleep 1
    /usr/bin/debug-app nvram write ssd
    sleep 1
    /usr/bin/debug-app ethernet nego disable
    sleep 1
    /usr/bin/debug-app bootmode check
    sleep 1
    /usr/bin/debug-app bootmode phy start
    sleep 1
    /usr/bin/debug-app stm32 bootcondition
    sleep 1
    /usr/bin/debug-app boot systemlog
    sleep 1
    /usr/bin/debug-app led set status green
    sleep 1
    /usr/bin/debug-app ethernet setip 2 192.168.10.7
    sleep 1
    /usr/bin/debug-app check all 2
    sleep 1
    /usr/bin/debug-app bootmode complete
    sleep 1
    /usr/bin/debug-app bootmode startSh
    ;;
  stop)
    echo "Stopping dhcpcd"
    killall dhcpcd
    ;;
  restart)
    echo "Restarting dhcpcd"
    killall dhcpcd
    /sbin/dhcpcd
    ;;
  *)
    echo "Usage: /etc/init.d/dhcpcd {start|stop|restart}"
    exit 1
    ;;
esac

exit 0
